# quadlets/netbird-client.container
[Unit]
Description=Netbird Client
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/netbirdio/netbird:latest
ContainerName=netbird-client
HostName=%H

# Capabilities needed for VPN tunnel
AddCapability=NET_ADMIN
AddCapability=SYS_RESOURCE
AddCapability=NET_RAW
AddDevice=/dev/net/tun

# Management server
Environment=NB_MANAGEMENT_URL=https://net.zehru.com
EnvironmentFile=/run/secrets/netbird-client.env

# Persist state across restarts
Volume=netbird-client.volume:/etc/netbird:Z

Network=host

[Service]
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target