# quadlets-rootless/caddy-proxy.container
# quadlets-rootless/caddy-proxy.container
# Rootless quadlet for Caddy reverse proxy — runs as core, not root.
# Deployed by gitops-sync to /home/core/.config/containers/systemd/

[Unit]
Description=Caddy Reverse Proxy (rootless)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/thatseanwalsh/caddy-custom:latest

Volume=/home/core/caddy/Caddyfile:/etc/caddy/Caddyfile:ro,Z
Volume=caddy-data.volume:/data:Z
Volume=caddy-config.volume:/config:Z
EnvironmentFile=/run/user/1000/secrets/cloudflare.env

# Expose HTTPS and HTTP to host
PublishPort=80:80
PublishPort=443:443
# QUIC / HTTP3
PublishPort=443:443/udp

Network=caddy-net.network
NoNewPrivileges=true

[Service]
Restart=always
RestartSec=10s

[Install]
WantedBy=default.target