# archive/zitadel.container
[Unit]
Description=Zitadel Identity Provider
After=network-online.target zitadel-db.service
Wants=network-online.target
Requires=zitadel-db.service caddy-net.network zitadel-db.volume zitadel-certs.volume

[Container]
Image=ghcr.io/zitadel/zitadel:latest
Network=caddy-net.network

EnvironmentFile=/run/user/1000/secrets/zitadel.env

Exec=start --tlsMode external --masterkeyFromEnv

NoNewPrivileges=true
DropCapability=ALL

HealthCmd=wget -qO- http://localhost:8080/debug/ready || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=180

[Install]
WantedBy=default.target