# quadlets-rootless/921/netbird-signal.container
# /etc/containers/systemd/netbird-signal.container
# Quadlet unit for the Netbird signal server

[Unit]
Description=Netbird Signal Server
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/netbirdio/signal:latest
Network=caddy-net.network
NoNewPrivileges=true
DropCapability=ALL
Exec=run --port 10000 --log-level info

[Service]

Restart=always
RestartSec=10s

[Install]
WantedBy=default.target
