# quadlets-rootless/921/zitadel-db.container
[Unit]
Description=PostgreSQL for Zitadel
After=network-online.target
Wants=network-online.target

[Container]
Image=docker.io/postgres:16-alpine
Network=caddy-net.network
Volume=zitadel-db.volume:/var/lib/postgresql/data:Z

Environment=POSTGRES_USER=zitadel
Environment=POSTGRES_DB=zitadel
EnvironmentFile=/run/user/1000/secrets/zitadel-db.env

NoNewPrivileges=true

HealthCmd=pg_isready -U zitadel
HealthInterval=10s
HealthTimeout=5s
HealthRetries=5

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=120

[Install]
WantedBy=default.target